# configs/git-syr.service
[Unit]
Description=git-syr
Documentation=https://github.com/aceberg/git-syr
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/git-syr -c /etc/git-syr/config.yaml -r /etc/git-syr/repos.yaml -l /var/log/git-syr.log
Restart=on-failure

[Install]
WantedBy=multi-user.target